# service.conf
start on startup
task
exec /home/Desktop/script.sh