# service.conf
#başlagıçta dosyayı bulup çalıştır

start on startup
task
exec /home/Desktop/script.sh
